# MonitoreoWeblogic-main/weblogic_monitor/deployment/weblogic-monitor.service
[Unit]
Description=WebLogic Monitor Backend API
After=network.target

[Service]
Type=simple
User=oracle
Group=oinstall
WorkingDirectory=/opt/weblogic-monitor/backend
Environment="PATH=/opt/weblogic-monitor/backend/venv/bin:/u01/oracle/middleware/oracle_common/common/bin"
Environment="ORACLE_HOME=/u01/oracle/middleware"
Environment="MW_HOME=/u01/oracle/middleware"
Environment="WLS_HOME=/u01/oracle/middleware/wlserver"
ExecStart=/opt/weblogic-monitor/backend/venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000 --workers 4
Restart=always
RestartSec=10

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=weblogic-monitor

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
